# target/linux/ar71xx/base-files/etc/uci-defaults/99_network-config
!/bin/sh

board=$(ar71xx_board_name)

case "$board" in
ng-nodedon)
	uci set wireless.@wifi-device[0].disabled=0
	uci commit wireless
	wifi
	;;
esac

exit 0
